# meta-quanta/meta-olympus-nuvoton/recipes-x86/chassis/x86-power-control/obmc-chassis-poweron.target
[Unit]
Description=Chassis Power On
Conflicts=obmc-chassis-poweroff.target
RefuseManualStart=yes
RefuseManualStop=yes
